# testing/systemd/units/default.target
[Unit]
Description=Kevlar Default Target
AllowIsolate=yes
Wants=kevlar-getty.service
