# testing/systemd/units/default.target
[Unit]
Description=Kevlar Default Target
Requires=basic.target
After=basic.target
AllowIsolate=yes
